# pkgs/configure-firefox/configure-firefox.service
[Unit]
Description=Configure Firefox prefs and extensions

[Service]
Type=oneshot
ExecStart=@CONFIGURE-FIREFOX@

[Install]
WantedBy=default.target
